# 04_CultiPi/01_Software/02_cultiConf/etc/cron.hourly/cultipi
#!/bin/bash

exec 2>&1
exec > >(logger -t "culticron.hourly")

echo "Updating repositories..."
apt-get update
echo "...ok"
