# trunk/01_software/01_install/05_cultipi/conf-package/DEBIAN-cultitime/postinst
#!/bin/bash

if [ -f /etc/init.d/cultitime ]; then
    update-rc.d cultitime defaults
fi

/etc/init.d/cultitime force-reload
