# 01_software/01_install/05_cultipi/conf-package/DEBIAN-culticam/postinst
#!/bin/bash

if [ -f /etc/init.d/culticam ]; then
    update-rc.d culticam defaults
fi

/etc/init.d/culticam force-reload
