# 01_software/01_install/05_cultipi/conf-package/DEBIAN-culticonf/postinst
#!/bin/bash


sed -i "/www-data ALL = (root) NOPASSWD: \/etc\/init.d\/networking/d" /etc/sudoers
sed -i "/www-data ALL = (root) NOPASSWD: \/usr\/sbin\/update-rc.d/d" /etc/sudoers

