# trunk/01_software/01_install/05_cultipi/conf-package/DEBIAN-culticam/postinst
#!/bin/bash

if [ -f /etc/init.d/culticam ]; then
    update-rc.d culticam defaults
fi

if [  "`cat /etc/sudoers|grep cultipi|grep fswebcam`" == "" ]; then
        echo "cultipi ALL = (root) NOPASSWD: /usr/bin/fswebcam" >> /etc/sudoers
fi  

chmod +x /etc/init.d/culticam

/etc/init.d/culticam force-reload
