# 04_CultiPi/01_Software/02_cultiConf/etc/cron.daily/cultipi
#!/bin/bash

if [ "`apt-get -u upgrade --assume-no|grep cultipi`" != "" ]; then
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" cultipi
elif [ "`dpkg -s cultipi 2>/dev/null|grep Version`" == "" ]; then
    apt-get install -y  --force-yes cultipi
fi

if [ "`apt-get -u upgrade --assume-no|grep cultibox`" != "" ]; then
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" cultibox
    rm -Rf /var/cache/lighttpd/compress/*
    if [ "$1" == "" ]; then
        /etc/init.d/lighttpd force-reload
    fi
elif [ "`dpkg -s cultibox 2>/dev/null|grep Version`" == "" ]; then
    apt-get install -y  --force-yes cultibox
fi

if [ "`apt-get -u upgrade --assume-no|grep cultiraz`" != "" ]; then
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" cultiraz
elif [ "`dpkg -s cultiraz 2>/dev/null|grep Version`" == "" ]; then
    apt-get install -y  --force-yes cultiraz
fi


if [ "`apt-get -u upgrade --assume-no|grep cultitime`" != "" ]; then
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" cultitime
elif [ "`dpkg -s cultitime 2>/dev/null|grep Version`" == "" ]; then
    apt-get install -y  --force-yes cultitime
fi

if [ "`apt-get -u upgrade --assume-no|grep culticonf`" != "" ]; then
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" culticonf
elif [ "`dpkg -s culticonf 2>/dev/null|grep Version`" == "" ]; then
    apt-get install -y  --force-yes culticonf
fi


if [ -d /var/www/cultibox/tmp/import ]; then
    rm -Rf /var/www/cultibox/tmp/import/*
fi

if [ -d /var/www/cultibox/tmp/export ]; then
    rm -Rf /var/www/cultibox/tmp/export/*
fi

