# 04_CultiPi/01_Software/02_cultiConf/etc/cron.daily/cultipi
#!/bin/bash

while [ $# -gt 0 ]
do
    case "$1" in
        (--now) now=yes;;
        (--manual) manual=yes;;
        (*)  break;;
    esac
    shift
done


if [ "$manual" != "yes" ]; then
    exec 2>&1
    exec > >(logger -t "culticron.daily")
fi

#Random sleep:
if [ "$now" != "yes" ]; then
    sleep $[ ( $RANDOM % 3600 )  + 1 ]s
fi


#If there was an error whith dpkg/apt:
/usr/bin/dpkg --configure -a --force-confdef --force-confold

echo "--> Updating repositories..."
apt-get update -o Dir::Etc::sourcelist="sources.list.d/cultibox.list"  -o Dir::Etc::sourceparts="-" -o APT::Get::List-Cleanup="0"

for file in "`ls /etc/apt/sources.list.d/cultibox-*.list 2>/dev/null`"; do
    if [ "$file" != "" ]; then
        apt-get update -o Dir::Etc::sourcelist="$file" -o Dir::Etc::sourceparts="-" -o APT::Get::List-Cleanup="0"
    fi
done
echo "...ok"

if [ "`apt-get -u upgrade --assume-no|grep cultipi`" != "" ]; then
    echo "--> Upgrading cultipi package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" cultipi
    echo "...ok"
elif [ "`dpkg -s cultipi 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing cultipi package..."
    apt-get install -y  --force-yes cultipi
    echo "...ok"
fi

if [ "`apt-get -u upgrade --assume-no|grep cultibox`" != "" ]; then
    echo "--> Upgrading cultibox package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" cultibox
    rm -Rf /var/cache/lighttpd/compress/*
    /etc/init.d/lighttpd force-reload
    echo "...ok"
elif [ "`dpkg -s cultibox 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing cultibox package..."
    apt-get install -y  --force-yes cultibox
    echo "...ok"
fi

if [ "`apt-get -u upgrade --assume-no|grep cultiraz`" != "" ]; then
    echo "--> Upgrading cultiraz package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" cultiraz
    echo "...ok"
elif [ "`dpkg -s cultiraz 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing cultiraz package..."
    apt-get install -y  --force-yes cultiraz
    echo "...ok"
fi


if [ "`apt-get -u upgrade --assume-no|grep cultitime`" != "" ]; then
    echo "--> Upgrading cultitime package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" cultitime
    echo "...ok"
elif [ "`dpkg -s cultitime 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing cultitime package..."
    apt-get install -y  --force-yes cultitime
    echo "...ok"
fi

if [ "`apt-get -u upgrade --assume-no|grep cultidoc`" != "" ]; then
    echo "--> Upgrading cultidoc package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" cultidoc
    echo "...ok"
elif [ "`dpkg -s cultidoc 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing cultidoc package..."
    apt-get install -y  --force-yes cultidoc
    echo "...ok"
fi

if [ "`apt-get -u upgrade --assume-no|grep culticam`" != "" ]; then
    echo "--> Upgrading culticam package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" culticam
    echo "...ok"
elif [ "`dpkg -s culticam 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing culticam package..."
    apt-get install -y  --force-yes culticam
    echo "...ok"
fi

if [ "`apt-get -u upgrade --assume-no|grep culticonf`" != "" ]; then
    echo "--> Upgrading culticonf package..."
    apt-get install -y  --force-yes --only-upgrade -o Dpkg::Options::="--force-confdef" -o Dpkg::Options::="--force-confold" -o Dpkg::Options::="--force-confmiss" culticonf
    echo "...ok"
elif [ "`dpkg -s culticonf 2>/dev/null|grep Version`" == "" ]; then
    echo "--> Installing culticonf package..."
    apt-get install -y  --force-yes culticonf
    echo "...ok"
fi



if [ -d /var/www/cultibox/tmp/import ]; then
    echo "--> Cleaning import directory..."
    rm -Rf /var/www/cultibox/tmp/import/*
    echo "...ok"
fi

if [ -d /var/www/cultibox/tmp/export ]; then
    echo "--> Cleaning cultibox directory..."
    rm -Rf /var/www/cultibox/tmp/export/*
    echo "...ok"
fi

if [ "`ls /var/www/cultibox/tmp/webcam*.jpg 2>/dev/null`" != "" ]; then
    echo "--> Cleaning webcam files..."
    rm /var/www/cultibox/tmp/webcam*.jpg
    echo "...ok"
fi

if [ -f /etc/culticonf/packages ]; then
    packages="`cat /etc/culticonf/packages`"
    if [ "$packages" != "" ]; then
        apt-get update
        for pack in $packages; do
            if [ "`dpkg -l $pack 2>/dev/null|grep ii`" == "" ]; then
                echo "--> Installing package $pack..."
                apt-get install -y --force-yes $pack
                echo "...ok"

                if [ -f /etc/culticonf/postinst/$pack ]; then
                    echo "--> Executing postinstall script for $pack"
                    bash /etc/culticonf/postinst/$pack
                    echo "...ok"
                fi
            fi
        done
        apt-get clean
    fi
fi

if [ -f /var/lock/REBOOT_NEEDED ]; then
    echo "--> Rebooting system by lock file..."
    rm -Rf /var/lock/REBOOT_NEEDED
    shutdown -r now
    echo "...ok"
fi

uptime=$(</proc/uptime)
uptime=${uptime%%.*}

if [ "$uptime" != "" ]; then
    if [ $uptime -gt 432000 ]; then
        echo "--> Rebooting system by uptime..."
        shutdown -r now
        echo "...ok"
    fi
fi


/sbin/iw reg set FR

if [ ! -f /var/lock/UPGRADE ]; then
    nohup /bin/bash -c 'sleep 3600 && /etc/cron.daily/cultipi' >/dev/null 2>&1 &
    touch /var/lock/UPGRADE
else
    rm -f /var/lock/UPGRADE
fi

echo "--> Cleaning tmp directory..."
rm -Rf /tmp/*
echo "...ok"
