# 01_software/01_install/05_cultipi/conf-package/DEBIAN-cultipi/postinst
#!/bin/bash

if [ -f /etc/init.d/cultipi ]; then
    update-rc.d cultipi defaults
    /etc/init.d/cultipi restart
fi

if [ ! -f /var/log/cultipi ]; then
    touch /var/log/cultipi
    chown cultipi:cultipi /var/log/cultipi
fi

if [ -f /etc/cron.daily/cultipi ]; then
    chmod +x /etc/cron.daily/cultipi
fi

if [  "`cat /etc/sudoers|grep www-data|grep ifup`" == "" ]; then
    echo "cultipi ALL = (root) NOPASSWD: /sbin/modprobe" >> /etc/sudoers
fi


